# conf/firewall.service
[Unit]
Description=Firewall IPTables

[Service]
ExecStart=/usr/bin/firewall

[Install]
WantedBy=multi-user.target